# files/.config/systemd/user/app-org.fcitx.Fcitx5@autostart.service.d/override.conf
[Service]
Restart=on-failure
RestartSec=2s

[Unit]
StartLimitIntervalSec=30s
StartLimitBurst=5
