# files/.config/systemd/user/acme-renew.timer
[Unit]
Description=Daily acme.sh renew (user)

[Timer]
OnCalendar=daily
RandomizedDelaySec=1h
Persistent=true

[Install]
WantedBy=timers.target
